# rootfs/etc/services.d/node-2/run
#!/bin/bash

echo "* Starting node 2..."
grin server -c /data/node-2/grin-server.toml run

s6-svc -d /var/run/s6/services/node-2
