# rootfs/etc/services.d/node-1/run
#!/bin/bash

echo "* Starting node 1..."
grin server -c /data/node-1/grin-server.toml

s6-svc -d /var/run/s6/services/node-1
