# rootfs/etc/services.d/wallet-2/run
#!/bin/bash

echo "* Starting wallet RPC 2..."
grin-wallet --usernet -p password -r http://127.0.0.1:3513 listen -l 3515
