# rootfs/etc/services.d/wallet-2/run
#!/bin/bash

echo "* Starting wallet RPC 2..."
grin-wallet -d /data/wallet-2 -p password listen
